# Debian/postinst
systemctl daemon-reload
systemctl enable sexy.service
deb-systemd-invoke restart sexy
